# deploy/launchd/com.finaq.cio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        FINAQ CIO heartbeat — fires twice a day on the user's Mac (PT).

        Cadence:  5am PT + 1pm PT  =  8am ET + 4pm ET (US market hours).
        Mode:     `auto` — dispatcher decides heartbeat vs catchup based on
                  freshness (>8h since last successful cycle = catchup).

        RunAtLoad=true:  fire one cycle when the launchd job is loaded.
        Combined with the freshness check in `cio.dispatcher._resolve_auto_mode`,
        this means: when the Mac wakes from sleep / boots / the user logs in,
        if a heartbeat slot was missed (lid was closed at 5am), we run a
        catchup cycle instead of waiting for the next slot.

        Install:  see `deploy/README.md` — copy to ~/Library/LaunchAgents/
                  and `launchctl load`. Edit `<key>WorkingDirectory</key>` and
                  the python path below to match your local setup before load.
    -->
    <key>Label</key>
    <string>com.finaq.cio</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use the venv's Python so the import path is correct. The venv
             must already be set up via `pip install -r requirements.txt`. -->
        <string>/Users/jujot/developer/finaq/.venv/bin/python</string>
        <string>-m</string>
        <string>cio.dispatcher</string>
        <string>--mode</string>
        <string>auto</string>
    </array>

    <!-- Critical: cd into the project root so .env is loaded by python-dotenv,
         and so relative paths (`theses/`, `data_cache/`) resolve correctly. -->
    <key>WorkingDirectory</key>
    <string>/Users/jujot/developer/finaq</string>

    <!-- Schedule: 5am + 1pm local time (PT for the user). launchd treats
         these as wall-clock schedules; if the Mac is asleep at the firing
         time, the job runs at next wake (RunAtLoad=true ensures a boot-path
         catch-up too). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Fire on load (boot, login, manual `launchctl load`). The dispatcher's
         freshness check picks the right mode (heartbeat vs catchup) so the
         user doesn't get pinged twice when they boot during a normal slot. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Single instance: don't fire a second cycle if the previous one is
         still running. Two heartbeats in flight would race on state.db
         and double-bill OpenRouter. -->
    <key>AbandonProcessGroup</key>
    <false/>

    <!-- Stdout / stderr → log files in data_cache/ so we can grep history
         without leaving the project tree. Rotate manually. -->
    <key>StandardOutPath</key>
    <string>/Users/jujot/developer/finaq/data_cache/cio.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jujot/developer/finaq/data_cache/cio.err.log</string>
</dict>
</plist>
